# 10-systemd/echod.socket
[Unit]
Description=Echo daemon activation socket

[Socket]
ListenStream=/run/isel/echod/request

[Install]
WantedBy=sockets.target
